# .config/systemd/user/batterymon.service
[Unit]
Description=Batterie überwachen

[Service]
Type=oneshot
ExecStart=/home/oli/.local/bin/batmonitor