# .config/systemd/user/synccalendar.timer
[Unit]
Description=Timer für Kalenderaktualisierung

[Timer]
OnBootSec=1min
OnUnitActiveSec=30min
Persistent=true

[Install]
WantedBy=timers.target
