# .config/systemd/user/synccalendar.service
[Unit]
Description=Kalender herunterladen

[Service]
Type=oneshot
ExecStart=/home/oli/.local/bin/synccalendar